# project/korra_common/packages/wdcore_avatar/lib/systemd/system/fstrim.service
[Unit]
Description=Discard unused blocks
After=primax-mdev.service

[Service]
Type=oneshot
ExecStart=/sbin/fstrim -a
