# project/korra_common/packages/wdcore_avatar/lib/systemd/system/multi-user.target.wants/primax-confchk.service
[Unit]
Description=Primax confchk RC Script (S99confchk)
DefaultDependencies=no
After=primax-rtc.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S99confchk start"
