# project/korra_common/packages/wdcore_avatar/lib/systemd/system/multi-user.target.wants/primax-netatalk.service
[Unit]
Description=Primax netatalk RC Script (S55netatalk)
DefaultDependencies=no
After=primax-ifplugd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S55netatalk start"
ExecStop=/bin/bash -c "/etc/init.d/S55netatalk stop"
