# project/korra_common/packages/wdcore_avatar/lib/systemd/system/multi-user.target.wants/primax-apacs.service
[Unit]
Description=Primax AP Auto Channel Selection RC Script (S99apacs)
DefaultDependencies=no
After=primax-confchk.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S99apacs start"
