# project/korra_common/packages/wdcore_avatar/lib/systemd/system/multi-user.target.wants/primax-multi-role.service
[Unit]
Description=Primax multi-role RC Script (S73multi-role)
DefaultDependencies=no
After=primax-ifplugd.service

[Service]
Type=forking
ExecStart=/bin/bash -c "/etc/init.d/S73multi-role start"
ExecStop=/bin/bash -c "/etc/init.d/S73multi-role stop"
