# project/korra_common/packages/wdcore_avatar/etc/cron.d/auto_fstrim
00 3 * * * root  /sbin/fstrim_trigger.sh
