# project/korra_common/packages/wdcore_avatar/etc/cron.d/auto_update
00 3 * * * root  /usr/local/sbin/getNewFirmwareUpgrade.sh immediate send_alert; /usr/local/sbin/checkAutoUpdate.sh
