# project/korra_common/packages/wdcore_avatar/lib/systemd/system/sysinit.target.wants/primax-initdisk.service
[Unit]
Description=Primax mdev initdisk Script (S43initdisk)
DefaultDependencies=no
After=primax-mdev.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S43initdisk start"
ExecStop=/bin/bash -c "/etc/init.d/S43initdisk stop"
