# project/korra_common/packages/wdcore_avatar/lib/systemd/system/sysinit.target.wants/primax-hostapd.service
[Unit]
Description=Primax hostapd RC Script (S60hostapd)
DefaultDependencies=no
After=primax-wilink8.service
Before=primax-mdev.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S60hostapd start"
ExecStop=/bin/bash -c "/etc/init.d/S60hostapd stop"
TimeoutSec=30s
