# project/korra_common/packages/wdcore_avatar/lib/systemd/system/sysinit.target.wants/primax-urandom.service
[Unit]
Description=Primax urandom RC Script (S15urandom)
DefaultDependencies=no
After=systemd-remount-fs.service
Before=sysinit.target shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S15urandom start"
ExecStop=/bin/bash -c "/etc/init.d/S15urandom stop"
