# project/korra_common/packages/wdcore_avatar/lib/systemd/system/sysinit.target.wants/primax-factory-restore.service
[Unit]
Description=Primax factory-restore RC Script (S47factory-restore)
DefaultDependencies=no
After=systemd-remount-fs.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "/etc/init.d/S47factory-restore start"
ExecStop=/bin/bash -c "/etc/init.d/S47factory-restore stop"
